# etc/Supervisor/laravel_vue_sample_app_worker.conf
[program:laravel_vue_sample_app_worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/laravel-vue-sample-app/artisan queue:work --queue=default --sleep=3 --tries=3 --backoff=3
numprocs=8
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
stopwaitsecs=3600
#redirect_stderr=true
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
stderr_logfile=/var/www/laravel-vue-sample-app/storage/logs/worker_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
